# web/crontabs/root
# DO NOT EDIT THIS FILE - edit the master and reinstall.
# (/tmp/crontab.2JXfY5/crontab installed on Sat Oct 30 19:38:39 2021)
# (Cron version -- $Id: crontab.c,v 2.13 1994/01/17 03:20:37 vixie Exp $)
#
# Environment Variables necessary to run the Scan Scheduler script
#
PATH=/scripts:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin:/app
* * * * * cd /app && /usr/bin/python3 /app/scheduled_tasks_runner.py

